# stress/stress.service
[Unit]
Description=stress
[Service]
ExecStart=/usr/bin/stress --cpu 6
[Install]
WantedBy=multi-user.target
