# open-webui/config/open-webui.service
[Unit]
Description=open-webui-backend
After=network.target
[Service]
Type=simple
EnvironmentFile=/root/open-webui/.env
ExecStart=/bin/bash -c 'source /root/anaconda3/etc/profile.d/conda.sh && conda activate open-webui && /root/open-webui/backend/start.sh'
Restart=always
RestartSec=10s
[Install]
WantedBy=multi-user.target
